# deploy/uwsgi.conf
# simple uWSGI script

description "uwsgi aprs2netportal"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --ini /var/www/aprs2netportal/deploy/uwsgi.ini
